# data_fetch/systemd/aviation-hub.service
[Unit]
Description=Aviation Hub Ingestion Service
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=craig
Group=craig
WorkingDirectory=/home/craig/projects/Aviation-hub/data_fetch
# Optional: Environment="VATSIM_BOOKINGS_API_KEY=..." or EnvironmentFile=-/etc/aviation-hub.env
ExecStart=/home/craig/projects/Aviation-hub/data_fetch/.venv/bin/python /home/craig/projects/Aviation-hub/data_fetch/src/main.py
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
