# aviation_hub/systemd/aviation-hub.service
[Unit]
Description=Aviation Hub Ingestion Service
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=CHANGE_ME
WorkingDirectory=/path/to/aviation_hub
ExecStart=/path/to/aviation_hub/.venv/bin/python /path/to/aviation_hub/src/main.py
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
